# docker/api/cron/crontab
* * * * * su - www-data /bin/bash -c "/usr/local/bin/php /var/www/artisan schedule:run >> /var/log/laravel-schedule-%Y%m%d.%H%M%S.log 2>&1"
